# root/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

crontab /etc/crontabs/root
exec \
    /usr/sbin/cron -f -L 15